# systemd/files-agent-node.service
[Unit]
Description=File Panel Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=filepanel
Group=filepanel
WorkingDirectory=/opt/files-agent
EnvironmentFile=/etc/files-agent/files-agent.env
ExecStart=/opt/files-agent/.venv/bin/python -m app.agent_main
Restart=on-failure
RestartSec=3
PrivateTmp=true
UMask=0027

[Install]
WantedBy=multi-user.target
